# wireguard-peer/root/etc/services.d/api/run
#!/usr/bin/with-contenv bash
_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down wg0
}

trap _term SIGTERM
python3 app.py & wait