# zigbee2mqtt-multi-coordinator/root/etc/services.d/zigbee2mqtt-multi-coordinator/run
#!/usr/bin/with-contenv bashio

export NODE_PATH=/opt/zigbee2mqtt/node_modules
COORDINATORS=$(bashio::config 'coordinators')
cd /opt/zigbee2mqtt

for COORDINATOR in $COORDINATORS
do
    NAME=$(echo $COORDINATOR | jq -r '.name')
    bashio::log.info "Handing over control to Zigbee2mqtt Core for coordinator $NAME..."
    ZIGBEE2MQTT_DATA=/config/zigbee2mqtt-multi-coordinator/$NAME npm start &
done
sleep infinity & wait
