# zigbee2mqtt-multi-coordinator/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

nginx
sleep infinity & wait